# config/forum_init.conf
description "Fórum"

start on runlevel [23]
stop on shutdown

exec sudo sh -c "bash /home/ubuntu/apps/forum/current/config/unicorn_init.sh start"

respawn